# scripts/rosev_ircsystem.service
[Unit]
Description=ReactOS Deutschland e.V. IRC System

[Service]
ExecStart=/usr/local/bin/rosev_ircsystem --daemon /srv/irc
Type=forking
